# salt/logging/files/cron.daily/00logwatch
#!/bin/bash

test -x /usr/share/logwatch/script/logwatch.pl || exit 0

/usr/sbin/logwatch --mailto segarock@gmail.com